# Bluetooth_Temp/crontab
* * * * *  /home/yaofei/orangepi_cfgs/Bluetooth_Temp/read_w1.py
